# CodexTokenMonitor/config/launch_agents/local.codex-usage-monitor.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.codex-usage-monitor.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
    <string>-m</string>
    <string>codex_usage_monitor.cli</string>
    <string>dashboard</string>
    <string>--port</string>
    <string>8766</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/linzezhang/Documents/Codex/2026-06-02/new-chat/outputs/codex-usage-monitor</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/linzezhang/Library/Logs/CodexUsageMonitor/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/linzezhang/Library/Logs/CodexUsageMonitor/launchd.err.log</string>
</dict>
</plist>
